# deploy/hypha-whisper-watchdog.service
[Unit]
Description=Hypha Whisper Node — health watchdog
Documentation=https://github.com/reef-imaging/hypha-whisper-node
After=network-online.target hypha-whisper.service
Wants=network-online.target
BindsTo=hypha-whisper.service
PartOf=hypha-whisper.service

[Service]
Type=simple
User=%USER%
WorkingDirectory=%WORKING_DIR%
Environment="PATH=%HOME%/.local/bin:/usr/local/bin:/usr/bin:/bin"

ExecStart=/usr/bin/python3 %WORKING_DIR%/watchdog.py

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hypha-whisper-watchdog

[Install]
WantedBy=multi-user.target
